# debian/vyatta-cfg.postrm.in
#!/bin/bash -e

if [ "$1" = "purge" ]; then
  update-rc.d vyatta-router remove >/dev/null || exit $?
  systemctl disable vyatta-router.service >/dev/null || exit $?
fi

