# debian/vyatta-cfg.postrm.in
#!/bin/bash -e

if [ "$1" = "purge" ]; then
  systemctl disable vyos-initialize.service >/dev/null || exit $?
  systemctl disable vyos-config-load.service >/dev/null || exit $?
fi

