# debian/vyatta-cfg.prerm.in
#!/bin/bash -e

if [ "$1" = "purge" ]; then
  systemctl disable vyos-router.service >/dev/null || exit $?
fi
