# debian/vyatta-cfg.postrm.in
#!/bin/bash

if [ "$1" = "purge" ]; then
  update-rc.d vyatta-ofr remove >/dev/null || exit $?
fi

