# debian/vyatta-cfg.postinst.in
#!/bin/bash

prefix=@prefix@
sysconfdir=@sysconfdir@

update-rc.d vyatta-ofr defaults 90 >/dev/null
# do we want to start vyatta-ofr here in postinst?

